# subtitle-product/jife-subtitles.service
[Unit]
Description=JIFE Japanese Subtitle System
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/vincent/JIFE/subtitle-product
Environment=AUDIO_DEVICE=plughw:2,0
# Create power control FIFO
ExecStartPre=/bin/bash -c 'rm -f /tmp/jife-power-control; mkfifo /tmp/jife-power-control; chmod 666 /tmp/jife-power-control'
ExecStartPre=/usr/bin/docker compose down --remove-orphans
# Start power control listener in background, then start docker
ExecStart=/bin/bash -c '/home/vincent/JIFE/subtitle-product/scripts/power-control.sh & /usr/bin/docker compose up'
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
